# x11/.config/autostart/keybindings.desktop
[Desktop Entry]
Name=Keybindings for Samus
Exec=/home/kofdog/.bin/samus-keybindings.sh
Terminal=true
Type=Application
